# etc/cron.d/openclaw-memory-ingest
# Description: Daily memory fact extraction from conversation history
0 8 * * * openclaw /usr/local/bin/memory-ingest-wrapper.sh 2>/dev/null
