# etc/cron.d/openclaw-gemini-rate
*/5 * * * * openclaw /home/openclaw/.openclaw/scripts/rate-limit-tracker.sh
