# ops/api-health-check.cron
* * * * * root /usr/local/bin/check-api-health.sh > /dev/null 2>&1
